# deploy_tools/gunicorn-systemd.template.service
[Unit]
Description=Gunicorn server for DOMAIN

[Service]
Restart=on-failure
User=grich
WorkingDirectory=/home/grich/sites/DOMAIN
EnvironmentFile=/home/grich/sites/DOMAIN/.env

ExecStart=/home/grich/sites/DOMAIN/virtualenv/bin/gunicorn \
    --bind unix:/tmp/DOMAIN.socket \
    superlists.wsgi:application

[Install]
WantedBy=multi-user.target